# Library/LaunchAgents/com.landonschropp.disable-screen-tint.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.landonschropp.disable-screen-tint</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/shortcuts</string>
      <string>run</string>
      <string>Disable Screen Tint</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/landon/Library/Logs/com.landonschropp.disable-screen-tint/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/landon/Library/Logs/com.landonschropp.disable-screen-tint/stderr.log</string>
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Hour</key>
        <integer>7</integer>
      </dict>
    </array>
  </dict>
</plist>
