# Library/LaunchAgents/com.landonschropp.update-packages.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<!--
For more information on configuring launchd, please read this article:
https://www.maketecheasier.com/use-launchd-run-scripts-on-schedule-macos/
-->

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.landonschropp.update-packages</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>exec $HOME/.dotfiles/bin/update-packages</string>
    </array>
    <key>RunAtLoad</key>
    <true />
    <key>StartInterval</key>
    <integer>14400</integer>
  </dict>
</plist>
